# files/etc/supervisor.d/sshd/sshd.ini
[program:sshd]
priority=100
command=/opt/local/bin/x-start_sshd.sh
user=root
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stdout_logfile_backups=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stderr_logfile_backups=0

